# scheduler/crontab
* * * * * curl http://backend:8000/run-coordinateur >> /tmp/cron.log 2>&1
