# scheduler/crontab
* 10 * * * curl http://backend:8000/dashboard-data >> /tmp/cron.log 2>&1
